# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000 --workers 3 --worker-class gevent
directory=/root/bbs
autostart=true
autorestart=true


